# roles/setup_node/templates/polkadot.service.j2
[Unit]
Description=Polkadot Validator

[Service]
ExecStart=/usr/local/bin/polkadot
Restart=always
RestartSec=120

[Install]
WantedBy=multi-user.target